# blog/environment/crontab
*/5 * * * * /usr/local/bin/php /var/www/html/bin/console app:subscription 2>&1 >> /var/log/app.log
